# tools/docker/cells/commands.conf
[program:cmd]

command=/bin/bash /home/cells/cmd.sh
directory=/home/cells

autorestart=true
autorestart=true
startretries=3

statsecs=10
